# procd.service
[Unit]
Description="Procd"
ConditionPathExists=/etc/procd.conf

[Service]
Type=exec
Restart=always
ExecStart=/usr/bin/procd
